# upsmon.service
[Unit]
Description = UPSmon monitoring script

[Service]
ExecStart = /usr/bin/python3 /opt/upsmon/upsmon.py

[Install]
WantedBy = multi-user.target